# supervisor/aprs.conf
[supervisord]
nodaemon=true

[program:aprs]
directory=/home/ssmith/micro-aprs
command="/home/ssmith/micro-aprs/start_igate.sh"
autostart=true
autorestart=true
redirect_stderr=true
